# package/glinet-privacy/files/etc/uci-defaults/99-glinet-privacy-device
#!/bin/sh
[ -x /usr/libexec/glinet-privacy/apply-device-profile.sh ] && \
	/usr/libexec/glinet-privacy/apply-device-profile.sh
exit 0
