# package/glinet-privacy/files/etc/uci-defaults/99-glinet-privacy-firewall
#!/bin/sh
[ -x /usr/bin/apply-privacy-firewall-includes.sh ] && /usr/bin/apply-privacy-firewall-includes.sh
exit 0
